# scripts/github-runner.service
[Unit]
Description=GitHub Actions Self-Hosted Runner
Documentation=https://github.com/nficano/github-runner
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=notify
User=github-runner
Group=github-runner
ExecStart=/usr/local/bin/github-runner start --config /etc/github-runner/config.toml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5s
WatchdogSec=60s
TimeoutStartSec=30s
TimeoutStopSec=60s

# Security hardening
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
NoNewPrivileges=yes
ReadWritePaths=/var/lib/github-runner /var/log/github-runner /var/run/github-runner
CapabilityBoundingSet=
AmbientCapabilities=
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallArchitectures=native

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Environment
Environment=HOME=/var/lib/github-runner
WorkingDirectory=/var/lib/github-runner

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=github-runner

[Install]
WantedBy=multi-user.target
